# solr/files/upstart.conf
{% set java = salt['pillar.get']('solr:java', {}) -%}
{% set initial_heap = java.initial_heap|d('128M') -%}
{% set maximum_heap = java.maximum_heap|d('256M') -%}
description     "Solr Service"
author          "Adam Venturella"
version         "0.0.1"

# Make sure the file system and network devices have started before
# we begin the daemon
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

kill timeout 30
respawn

script
    chdir {{ install_location }}
    # -Xms set initial Java heap size
    # -Xmx set maximum Java heap size
    exec /usr/bin/java -Xms{{ initial_heap }} -Xmx{{ maximum_heap }} -jar start.jar
end script
