# supervisor/files/program.conf
{% from "supervisor/files/macros.jinja" import process_value %}

[program:{{ name }}]
{% for key, value in values.iteritems() -%}
{{ key }}={{ process_value(value) }}
{% endfor -%}
